# service/update-cloudflare-dns.timer
[Unit]
Description=Update CloudFlare DNS entries every five minutes

[Timer]
OnBootSec=1min
OnUnitActiveSec=5min

[Install]
WantedBy=timers.target
